# files/startup/pencil.conf
# pencil

start on runlevel [345]
stop on starting shutdown

respawn
chdir /opt/pencil
exec /usr/lib/ruby/gems/1.8/gems/pencil-0.2.9/bin/pencil --config-dir=/opt/pencil/config
