# files/startup/pencil.conf
# pencil

start on runlevel [345]
stop on shutdown

respawn
chdir /opt/pencil
exec /usr/bin/pencil --config-dir=/opt/pencil/config
