# files/startup/sentry.conf
# sentry

start on runlevel [345]
stop on starting shutdown

respawn
chdir /opt/sentry
exec /usr/bin/sentry --config=sentry.conf.py start 
